# scripts/init/rc/basenine
#!/sbin/openrc-run

# Move this file into /etc/init.d/basenine
# To add Basenine service at boot time, run:
#   rc-update add basenine boot

command=/usr/local/bin/basenine
